# deploy/com.alphalete.sales-boards.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.alphalete.sales-boards</string>

    <!-- Path is the laptop placeholder; on the mini, install_agent.py rewrites it
         to the mini's repo path (python -m automations.day_orchestrator.install_agent
         sales-boards). -->
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/megan/1st Claude Folder/deploy/sales_boards.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/megan/1st Claude Folder</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTHONPATH</key>
        <string>/Users/megan/1st Claude Folder</string>
    </dict>

    <!-- DAILY incl. weekends (no Weekday key). First pass 5:10 CT (Carlos, via
         Megan 2026-07-23) — ten minutes after com.alphalete.vantura-slack-sales'
         05:00 pass writes the previous day's Base/BOX/AT&T cells from the
         channel. HISTORY: 6:00 from 2026-07-18 (matched the VA's ~5:58-6:01
         window), moved to 7:00 by Eve 2026-07-21 because it fired before the day
         cells were HAND-entered. Automating that fill is what makes 5:10 safe —
         so if the slack-sales agent is ever disabled, this MUST go back to 7:00
         or the images post a blank day. Then q25m retries; the module HOLDS on a
         wrong-week board and skips boards already in today's thread, so later
         passes are no-ops. The whole ladder moves together. TIME KNOB: edit
         these entries (Hour/Minute are machine LOCAL time; Lucy 2 is Central). -->
    <key>StartCalendarInterval</key>
    <array>
        <dict><key>Hour</key><integer>5</integer><key>Minute</key><integer>10</integer></dict>
        <dict><key>Hour</key><integer>5</integer><key>Minute</key><integer>35</integer></dict>
        <dict><key>Hour</key><integer>6</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Hour</key><integer>6</integer><key>Minute</key><integer>25</integer></dict>
        <dict><key>Hour</key><integer>6</integer><key>Minute</key><integer>50</integer></dict>
        <dict><key>Hour</key><integer>7</integer><key>Minute</key><integer>15</integer></dict>
        <dict><key>Hour</key><integer>7</integer><key>Minute</key><integer>40</integer></dict>
        <dict><key>Hour</key><integer>8</integer><key>Minute</key><integer>5</integer></dict>
    </array>

    <key>StandardOutPath</key>
    <string>/Users/megan/1st Claude Folder/output/logs/sales-boards.launchd.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/megan/1st Claude Folder/output/logs/sales-boards.launchd.err.log</string>
</dict>
</plist>
