# deploy/com.alphalete.box-order-log.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.alphalete.box-order-log</string>

    <!-- Path is the laptop placeholder; on the mini, install_agent.py rewrites it
         to that machine's repo path (python -m automations.day_orchestrator.install_agent
         box-order-log). -->
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/megan/1st Claude Folder/deploy/box_order_log.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/megan/1st Claude Folder</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTHONPATH</key>
        <string>/Users/megan/1st Claude Folder</string>
    </dict>

    <!-- DAILY incl. weekends (no Weekday key). TWICE a day, 7:00 and 8:30 CST —
         exactly what Carlos asked for (2026-07-18). 7:00 sits after the 6:00
         Sales Boards run so the two don't contend for the ownerville session,
         and is late enough that the Box extract has refreshed.

         ONLY ONE of these posts to Slack. box_order_log.sh writes a dated
         marker after a successful post, so 8:30 sees it and refreshes the
         sheet only. If 7:00 dies before posting there's no marker, and 8:30
         posts instead — a late post beats none, and it can never be two.

         TIME KNOB: edit these entries (Hour/Minute are machine LOCAL time;
         the mini is Central). -->
    <key>StartCalendarInterval</key>
    <array>
        <dict><key>Hour</key><integer>7</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Hour</key><integer>8</integer><key>Minute</key><integer>30</integer></dict>
    </array>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/tmp/com.alphalete.box-order-log.out</string>
    <key>StandardErrorPath</key>
    <string>/tmp/com.alphalete.box-order-log.err</string>
</dict>
</plist>
