# deploy/com.alphalete.new-start-followup-sat-pm.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.alphalete.new-start-followup-sat-pm</string>

    <!-- Path is the laptop placeholder; on the mini, regenerate with the install
         path via the python-replace step in workflows/scheduler-mini.md. -->
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/megan/1st Claude Folder/deploy/new_start_followup.sh</string>
        <string>--mode</string>
        <string>nudge</string>
        <string>--when</string>
        <string>auto</string>
        <string>--live</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/megan/1st Claude Folder</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTHONPATH</key>
        <string>/Users/megan/1st Claude Folder</string>
    </dict>

    <!-- Saturday 5:00pm Central (Weekday 6) — the last call as the field day
         ends. Split out of the -sat plist ON PURPOSE so NEITHER job exceeds 2
         calendar intervals: schedule_guard skips jobs with more than 2 as
         high-frequency pollers (schedule_guard.py `_timed_schedule`), which
         would drop them from the nightly anti-drift reload. Same wrapper and
         same flags as -sat; the clock picks the "last call" wording. -->
    <key>StartCalendarInterval</key>
    <array>
        <dict>
            <key>Weekday</key><integer>6</integer>
            <key>Hour</key><integer>17</integer>
            <key>Minute</key><integer>0</integer>
        </dict>
    </array>

    <key>StandardOutPath</key>
    <string>/Users/megan/1st Claude Folder/output/logs/new-start-followup-sat-pm.launchd.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/megan/1st Claude Folder/output/logs/new-start-followup-sat-pm.launchd.err.log</string>
</dict>
</plist>
