# deploy/com.alphalete.vantura-churn-daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.alphalete.vantura-churn-daily</string>

    <!-- Path is the laptop placeholder; on Lucy 2, install_agent.py rewrites it
         to that machine's repo path (lucy rerun install_vantura_churn_agent). -->
    <!-- --no-post added 2026-07-20: vantura keeps REFRESHING the board (full
         pull + write of LUCY CHURN / Activations) but no longer POSTS the churn
         screenshots — the consolidated b2b_metrics run (7:45am) posts them now,
         in order, with everything else. Revert this arg to restore posting. -->
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/megan/1st Claude Folder/deploy/vantura_churn_daily.sh</string>
        <string>--no-post</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/megan/1st Claude Folder</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTHONPATH</key>
        <string>/Users/megan/1st Claude Folder</string>
    </dict>

    <!-- Every day 7:00am, machine LOCAL time (Lucy 2 is Central). No Weekday
         key = fires all 7 days, matching the daily board refresh. -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>7</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/megan/1st Claude Folder/output/logs/vantura-churn-daily.launchd.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/megan/1st Claude Folder/output/logs/vantura-churn-daily.launchd.err.log</string>
</dict>
</plist>
