# deploy/com.alphalete.applicant-am.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.alphalete.applicant-am</string>

    <!-- Path is the laptop placeholder; on Lucy 1, install_agent.py rewrites it
         to that machine's repo path:
           python -m automations.day_orchestrator.install_agent applicant-am
         RUNS ON LUCY 1 — that machine holds the one-time ApplicantStream browser
         login (automations/applicant_tracker/.browser_profile). Francia's
         ApplicantStream → Applicant Tracker morning sync. -->
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/megan/1st Claude Folder/deploy/applicant_tracker.sh</string>
        <string>am</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/megan/1st Claude Folder</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTHONPATH</key>
        <string>/Users/megan/1st Claude Folder</string>
    </dict>

    <!-- 7:00am CST, Mon–Sat (Sun OFF). Reads YESTERDAY, so morning is fine.
         Runs export_call_list then update_second_round SEQUENTIALLY (one shared
         browser profile). launchd Weekday: 1=Mon … 6=Sat. TIME KNOB: edit the
         Hour/Minute below (machine LOCAL time; the mini is Central). -->
    <key>StartCalendarInterval</key>
    <array>
        <dict><key>Weekday</key><integer>1</integer><key>Hour</key><integer>7</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Weekday</key><integer>2</integer><key>Hour</key><integer>7</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Weekday</key><integer>3</integer><key>Hour</key><integer>7</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Weekday</key><integer>4</integer><key>Hour</key><integer>7</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Weekday</key><integer>5</integer><key>Hour</key><integer>7</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Weekday</key><integer>6</integer><key>Hour</key><integer>7</integer><key>Minute</key><integer>0</integer></dict>
    </array>

    <key>RunAtLoad</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/tmp/com.alphalete.applicant-am.out</string>
    <key>StandardErrorPath</key>
    <string>/tmp/com.alphalete.applicant-am.err</string>
</dict>
</plist>
