# deploy/com.alphalete.sales-boards.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.alphalete.sales-boards</string>

    <!-- Path is the laptop placeholder; on the mini, install_agent.py rewrites it
         to the mini's repo path (python -m automations.day_orchestrator.install_agent
         sales-boards). -->
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/megan/1st Claude Folder/deploy/sales_boards.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/megan/1st Claude Folder</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTHONPATH</key>
        <string>/Users/megan/1st Claude Folder</string>
    </dict>

    <!-- DAILY incl. weekends (no Weekday key). First pass 7:00 CT (Eve 2026-07-21;
         was 6:00 from 2026-07-18, which matched the VA's ~5:58-6:01 window but
         fired before the day cells are hand-entered), then q25m retries. The
         module HOLDS on a wrong-week board and skips boards already in today's
         thread, so later passes are no-ops. The whole ladder moves together —
         leaving an earlier pass in place would post before 7:00 and defeat the
         change. TIME KNOB: edit these entries (Hour/Minute are machine LOCAL
         time; Lucy 2 is Central). -->
    <key>StartCalendarInterval</key>
    <array>
        <dict><key>Hour</key><integer>7</integer><key>Minute</key><integer>0</integer></dict>
        <dict><key>Hour</key><integer>7</integer><key>Minute</key><integer>25</integer></dict>
        <dict><key>Hour</key><integer>7</integer><key>Minute</key><integer>50</integer></dict>
        <dict><key>Hour</key><integer>8</integer><key>Minute</key><integer>15</integer></dict>
        <dict><key>Hour</key><integer>8</integer><key>Minute</key><integer>40</integer></dict>
        <dict><key>Hour</key><integer>9</integer><key>Minute</key><integer>5</integer></dict>
        <dict><key>Hour</key><integer>9</integer><key>Minute</key><integer>30</integer></dict>
        <dict><key>Hour</key><integer>9</integer><key>Minute</key><integer>55</integer></dict>
    </array>

    <key>StandardOutPath</key>
    <string>/Users/megan/1st Claude Folder/output/logs/sales-boards.launchd.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/megan/1st Claude Folder/output/logs/sales-boards.launchd.err.log</string>
</dict>
</plist>
