# deploy/com.alphalete.carlos-captainship-headcount-mon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.alphalete.carlos-captainship-headcount-mon</string>

    <!-- Path is the laptop placeholder; on Lucy 2, install_agent.py rewrites it
         to that machine's repo path (lucy rerun install_carlos-captainship-headcount-mon_agent). -->
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/megan/1st Claude Folder/deploy/carlos_captainship_headcount_mon.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/megan/1st Claude Folder</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTHONPATH</key>
        <string>/Users/megan/1st Claude Folder</string>
    </dict>

    <!-- Monday 7:20am, machine LOCAL time (Lucy 2 is Central). launchd
         Weekday: 0/7=Sun, 1=Mon … 6=Sat.

         WHY 7:20 AND NOT 7:00: three browser reports on Lucy 2 fired at 7:00
         and they all share one Chrome profile (automations/uploaded/
         .browser_profile) — BOX Order Log, Vantura Churn, and this. Only one
         can hold it; the losers get "Opening in existing browser session" and
         die after a 4x8s retry budget. That killed this report on 2026-07-20:
         BOX took the profile 07:00:00-07:05:50, this started 07:00:05, waited
         32s and gave up. chrome_guard does NOT help — it protects automation
         Chrome by design and only closes a human's browser.
         7:20 clears BOX (~07:06) and Vantura (~07:12). Sales Boards runs 7:16
         but uses no browser, so it does not contend.
         If BOX or Vantura ever get slower, move this later — don't shorten
         theirs. -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>1</integer>
        <key>Hour</key>
        <integer>7</integer>
        <key>Minute</key>
        <integer>20</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/megan/1st Claude Folder/output/logs/carlos-captainship-headcount-mon.launchd.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/megan/1st Claude Folder/output/logs/carlos-captainship-headcount-mon.launchd.err.log</string>
</dict>
</plist>
