# deploy/com.alphalete.carlos-captainship-headcount-mon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.alphalete.carlos-captainship-headcount-mon</string>

    <!-- Path is the laptop placeholder; on Lucy 2, install_agent.py rewrites it
         to that machine's repo path (lucy rerun install_carlos-captainship-headcount-mon_agent). -->
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/megan/1st Claude Folder/deploy/carlos_captainship_headcount_mon.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/megan/1st Claude Folder</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTHONPATH</key>
        <string>/Users/megan/1st Claude Folder</string>
    </dict>

    <!-- Monday 7:00am, machine LOCAL time (Lucy 2 is Central). launchd
         Weekday: 0/7=Sun, 1=Mon … 6=Sat. -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>1</integer>
        <key>Hour</key>
        <integer>7</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/megan/1st Claude Folder/output/logs/carlos-captainship-headcount-mon.launchd.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/megan/1st Claude Folder/output/logs/carlos-captainship-headcount-mon.launchd.err.log</string>
</dict>
</plist>
