# deploy/com.alphalete.stf-field-check-11pm.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.alphalete.stf-field-check-11pm</string>

    <!-- Path is the laptop placeholder; on the mini it is rewritten to the
         install path by automations.day_orchestrator.install_agent (run via
         `lucy rerun install_stf_field_check_agent`). -->
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/megan/1st Claude Folder/deploy/stf_field_check_11pm.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/megan/1st Claude Folder</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTHONPATH</key>
        <string>/Users/megan/1st Claude Folder</string>
    </dict>

    <!-- Daily 11:00pm CST (23:00), machine LOCAL time (the mini is Central).
         Late enough that reps are out of the field and the day's knocks are
         final, but before the 4am board post. -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>23</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/megan/1st Claude Folder/output/logs/stf-field-check-11pm.launchd.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/megan/1st Claude Folder/output/logs/stf-field-check-11pm.launchd.err.log</string>
</dict>
</plist>
