# deploy/com.alphalete.orchestrator-schedule-guard.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.alphalete.orchestrator-schedule-guard</string>

    <!-- The committed laptop path; install_agent rewrites it to THIS machine's
         repo path on bootstrap (same as every other agent here). -->
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/megan/1st Claude Folder/deploy/orchestrator_schedule_guard.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/megan/1st Claude Folder</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTHONPATH</key>
        <string>/Users/megan/1st Claude Folder</string>
    </dict>

    <!-- 3:40am every day, in the machine's LOCAL time (mini is CST) — 20 minutes
         before the 4:00am day-orchestrator, so launchd re-loads the Hour-4
         schedule fresh from the plist file BEFORE it needs to fire. This is the
         self-heal for the recurring 4am→6am drift. -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>3</integer>
        <key>Minute</key>
        <integer>40</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/megan/1st Claude Folder/output/logs/orchestrator-schedule-guard.launchd.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/megan/1st Claude Folder/output/logs/orchestrator-schedule-guard.launchd.err.log</string>
</dict>
</plist>
