# deploy/com.alphalete.rc-autoread.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.alphalete.rc-autoread</string>

    <!-- Path is the laptop's; on the mini, regenerate with the install path via
         the python-replace step in workflows/scheduler-mini.md. -->
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/megan/1st Claude Folder/deploy/rc_autoread_10min.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/megan/1st Claude Folder</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTHONPATH</key>
        <string>/Users/megan/1st Claude Folder</string>
    </dict>

    <!-- Every 10 minutes on the :00/:10/:20/:30/:40/:50, 24/7. The wrapper
         (rc_autoread_10min.sh) gates the ACTIVE window to 7am-midnight CST;
         outside that it exits immediately. Mini is Central. -->
    <key>StartCalendarInterval</key>
    <array>
        <dict><key>Minute</key><integer>0</integer></dict>
        <dict><key>Minute</key><integer>10</integer></dict>
        <dict><key>Minute</key><integer>20</integer></dict>
        <dict><key>Minute</key><integer>30</integer></dict>
        <dict><key>Minute</key><integer>40</integer></dict>
        <dict><key>Minute</key><integer>50</integer></dict>
    </array>

    <!-- Don't stack runs if one is somehow still going (it finishes in seconds). -->
    <key>StandardOutPath</key>
    <string>/Users/megan/1st Claude Folder/output/logs/rc-autoread.launchd.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/megan/1st Claude Folder/output/logs/rc-autoread.launchd.err.log</string>
</dict>
</plist>
