# deploy/com.alphalete.bg-check-sync.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.alphalete.bg-check-sync</string>

    <!-- Path is the laptop placeholder; on the mini, regenerate with the install
         path via the python-replace step in workflows/scheduler-mini.md. -->
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/megan/1st Claude Folder/deploy/bg_check_sync.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/megan/1st Claude Folder</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTHONPATH</key>
        <string>/Users/megan/1st Claude Folder</string>
    </dict>

    <!-- 11:30 / 16:00 every day, machine LOCAL time (mini is Central).
         TWO entries (not 3) on purpose: the nightly schedule_guard only self-heals
         timed jobs with <=2 StartCalendarInterval entries (it treats 3+ as
         high-frequency pollers and skips them). A 3-entry version fell fully off
         launchd 2026-07-20 and nothing brought it back, so keep it at 2 to stay
         under the guard. 11:30 = the Monday must-run when the week's thread
         starts; 16:00 = the afternoon run that does the Friday repost (hour>=12). -->
    <key>StartCalendarInterval</key>
    <array>
        <dict>
            <key>Hour</key><integer>11</integer>
            <key>Minute</key><integer>30</integer>
        </dict>
        <dict>
            <key>Hour</key><integer>16</integer>
            <key>Minute</key><integer>0</integer>
        </dict>
    </array>

    <key>StandardOutPath</key>
    <string>/Users/megan/1st Claude Folder/output/logs/bg-check-sync.launchd.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/megan/1st Claude Folder/output/logs/bg-check-sync.launchd.err.log</string>
</dict>
</plist>
