# deploy/com.alphalete.social-scanner.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.alphalete.social-scanner</string>

    <!-- Path is the laptop's; on the mini, regenerate with the install path via
         `lucy rerun install_social_scanner_agent` (rewrites to the mini repo). -->
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/megan/1st Claude Folder/deploy/social_scanner_noon.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/megan/1st Claude Folder</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTHONPATH</key>
        <string>/Users/megan/1st Claude Folder</string>
    </dict>

    <!-- 12:00 noon every day, machine LOCAL time (the mini is Central). Single
         daily run — independent of the separate com.alphalete.brand-audit-noon
         job that also fires at noon. -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>12</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/megan/1st Claude Folder/output/logs/social-scanner.launchd.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/megan/1st Claude Folder/output/logs/social-scanner.launchd.err.log</string>
</dict>
</plist>
