# deploy/com.alphalete.machine-digest.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.alphalete.machine-digest</string>

    <!-- Path is the laptop's; on the mini, regenerate with the install path via
         `lucy install_machine_digest`. -->
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/megan/1st Claude Folder/deploy/machine_digest_daily.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/megan/1st Claude Folder</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTHONPATH</key>
        <string>/Users/megan/1st Claude Folder</string>
    </dict>

    <!-- Noon CST daily — matches Lucy 1's summary timing (the orchestrator's
         backstop_time is 12:00). The mini is Central, so Hour 12 = 12pm CST. -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key><integer>12</integer>
        <key>Minute</key><integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/megan/1st Claude Folder/output/logs/machine-digest.launchd.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/megan/1st Claude Folder/output/logs/machine-digest.launchd.err.log</string>
</dict>
</plist>
