# deploy/com.alphalete.new-start-followup-sun.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.alphalete.new-start-followup-sun</string>

    <!-- Path is the laptop placeholder; on the mini, regenerate with the install
         path via the python-replace step in workflows/scheduler-mini.md. -->
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/megan/1st Claude Folder/deploy/new_start_followup.sh</string>
        <string>--mode</string>
        <string>checklist</string>
        <string>--live</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/megan/1st Claude Folder</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTHONPATH</key>
        <string>/Users/megan/1st Claude Folder</string>
    </dict>

    <!-- Sunday 1:00pm Central (Weekday 0) — the slot Raf builds the list by hand
         today. Posts the numbered ✅ roll-up and tags whoever still hasn't sent. -->
    <key>StartCalendarInterval</key>
    <array>
        <dict>
            <key>Weekday</key><integer>0</integer>
            <key>Hour</key><integer>13</integer>
            <key>Minute</key><integer>0</integer>
        </dict>
    </array>

    <key>StandardOutPath</key>
    <string>/Users/megan/1st Claude Folder/output/logs/new-start-followup-sun.launchd.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/megan/1st Claude Folder/output/logs/new-start-followup-sun.launchd.err.log</string>
</dict>
</plist>
