# deploy/com.alphalete.lucy2-digest.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.alphalete.lucy2-digest</string>

    <!-- Runs on LUCY 1 (not Lucy 2) — it summarizes Lucy 2's runs from the
         shared Hub Activity log. Path is the laptop's; on the mini, regenerate
         via `lucy install_lucy2_digest`. -->
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/megan/1st Claude Folder/deploy/lucy2_digest_daily.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/megan/1st Claude Folder</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTHONPATH</key>
        <string>/Users/megan/1st Claude Folder</string>
    </dict>

    <!-- Noon CST daily — matches Lucy 1's summary timing. The mini is Central,
         so Hour 12 = 12pm CST. -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key><integer>12</integer>
        <key>Minute</key><integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/megan/1st Claude Folder/output/logs/lucy2-digest.launchd.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/megan/1st Claude Folder/output/logs/lucy2-digest.launchd.err.log</string>
</dict>
</plist>
