# deploy/com.alphalete.vantura-payroll-wed.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.alphalete.vantura-payroll-wed</string>

    <!-- Path is the laptop placeholder; on Lucy 2, install_agent.py rewrites it
         to that machine's repo path (lucy rerun install_vantura_payroll_agent). -->
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/megan/1st Claude Folder/deploy/vantura_payroll_wed.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/megan/1st Claude Folder</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTHONPATH</key>
        <string>/Users/megan/1st Claude Folder</string>
    </dict>

    <!-- Wednesday 11:00am, machine LOCAL time (Lucy 2 is Central). launchd
         Weekday: 0/7=Sun, 1=Mon, 2=Tue, 3=Wed, … 6=Sat. -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>3</integer>
        <key>Hour</key>
        <integer>11</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/megan/1st Claude Folder/output/logs/vantura-payroll-wed.launchd.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/megan/1st Claude Folder/output/logs/vantura-payroll-wed.launchd.err.log</string>
</dict>
</plist>
