# deploy/com.alphalete.att-churn-daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.alphalete.att-churn-daily</string>

    <!-- Path is the laptop placeholder; on Lucy 2, install_agent.py rewrites it
         to that machine's repo path (lucy rerun install_att_churn_agent). -->
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/megan/1st Claude Folder/deploy/att_churn_daily.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/megan/1st Claude Folder</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTHONPATH</key>
        <string>/Users/megan/1st Claude Folder</string>
    </dict>

    <!-- Every day 7:15am, machine LOCAL time (Lucy 2 is Central). AFTER vantura's
         7:00 refresh: the wrapper also hard-waits for any vantura_churn run to
         finish before touching Chrome, because both share CDP debug port 9246. -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>7</integer>
        <key>Minute</key>
        <integer>15</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/megan/1st Claude Folder/output/logs/att-churn-daily.launchd.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/megan/1st Claude Folder/output/logs/att-churn-daily.launchd.err.log</string>
</dict>
</plist>
