# deploy/com.alphalete.sara-down.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.alphalete.sara-down</string>

    <!-- Path is the laptop's; on the mini, regenerate with the install path via
         `lucy rerun install_sara_down_agent` (rewrites to the mini repo). -->
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/megan/1st Claude Folder/deploy/sara_down_5min.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/megan/1st Claude Folder</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTHONPATH</key>
        <string>/Users/megan/1st Claude Folder</string>
    </dict>

    <!-- Every 5 minutes, around the clock. SaraPlus can go down at 2am and the
         whole field is stuck, so no daytime-only gate here. -->
    <key>StartInterval</key>
    <integer>300</integer>

    <key>StandardOutPath</key>
    <string>/Users/megan/1st Claude Folder/output/logs/sara-down.launchd.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/megan/1st Claude Folder/output/logs/sara-down.launchd.err.log</string>
</dict>
</plist>
