# deploy/com.alphalete.resume-pushing.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.alphalete.resume-pushing</string>

    <!-- Path is the laptop placeholder; on Lucy 2, install_agent.py rewrites it
         to that machine's repo path (lucy rerun install_resume_pushing_agent). -->
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/megan/1st Claude Folder/deploy/resume_pushing_10min.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/megan/1st Claude Folder</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTHONPATH</key>
        <string>/Users/megan/1st Claude Folder</string>
    </dict>

    <!-- Fire every 10 min around the clock (:00/:10/:20/:30/:40/:50). The
         wrapper (resume_pushing_10min.sh) gates the ACTIVE window to 8am–10pm
         CST and SKIPS Saturday; outside that it exits immediately. Machine LOCAL
         time (Lucy 2 is Central). TIME KNOB for hours/days: edit the window +
         day gates in the wrapper, NOT this array. -->
    <key>StartCalendarInterval</key>
    <array>
        <dict><key>Minute</key><integer>0</integer></dict>
        <dict><key>Minute</key><integer>10</integer></dict>
        <dict><key>Minute</key><integer>20</integer></dict>
        <dict><key>Minute</key><integer>30</integer></dict>
        <dict><key>Minute</key><integer>40</integer></dict>
        <dict><key>Minute</key><integer>50</integer></dict>
    </array>

    <key>StandardOutPath</key>
    <string>/Users/megan/1st Claude Folder/output/logs/resume-pushing.launchd.out.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/megan/1st Claude Folder/output/logs/resume-pushing.launchd.err.log</string>
</dict>
</plist>
